# supervisord.conf
[supervisord]
nodaemon=true

[program:gateway-service]
command=/usr/local/go/bin/go run /app/main.go
autostart=true
autorestart=true
